# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:ams]
command=/etc/init.d/ams start
